# src/garagepi/systemd/plate_watcher.service
[Unit]
Description=GaragePi Plate Watcher (RTSP -> ALPR -> Garage)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/garagepi      # <-- change if your repo path differs
EnvironmentFile=/etc/default/plate_watcher
ExecStart=/usr/bin/python3 /home/pi/garagepi/automation/plate_watcher.py
Restart=always
RestartSec=5

# Hardening (optional but nice)
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
